# tools/leyes-senado/co.ricardoruiz.leyes-diario.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.ricardoruiz.leyes-diario</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ricardoruiz/ricardoruiz.co/tools/leyes-senado/run_diario.sh</string>
    </array>

    <!-- dos corridas: la lista del Senado se mueve varias veces al día
         (mañana y tarde). Baja solo lo nuevo, sin reactivar el WAF. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <!-- si la laptop estaba dormida a la hora agendada, launchd lo corre al despertar -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ricardoruiz/ricardoruiz.co/Bases de datos/leyes-senado/diario/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ricardoruiz/ricardoruiz.co/Bases de datos/leyes-senado/diario/launchd.err.log</string>
</dict>
</plist>
